# frontend/etc/sv/httpbuf/run
#!/bin/sh

export GOATCOUNTER_BUFFER_SECRET=$%httpbuf%$
exec chpst -u martin:martin /usr/local/bin/goatcounter buffer -backend http://139.162.3.42:8081
