# frontend/etc/sv/httpbuf/run
#!/bin/sh

exec chpst -u goatcounter:goatcounter /home/goatcounter/bin/httpbuf
