# frontend/etc/sv/varnishncsa/run
#!/bin/sh

exec /usr/bin/varnishncsa \
	-a -w /var/log/varnish/varnishncsa.log \
	-q 'ReqURL ne "/status"' \
	-F '%h "%{X-Forwarded-For}i" %t "%r" %s "%{Referer}i" "%{User-agent}i"'
