# frontend/etc/sv/hitch/run
#!/bin/sh

mkdir -p /var/lib/hitch
chown hitch:hitch /var/lib/hitch

ulimit -Hn 16384
ulimit -Sn 4096

exec /usr/sbin/hitch --config=/etc/hitch/hitch.conf /etc/hitch/pem/*.pem
