# shared/etc/sv/haveged/run
#!/bin/sh
exe 1>&2
exec /usr/sbin/haveged --Foreground
