# goatcounter/etc/periodic/weekly/data
#!/bin/sh

set -euC

su postgres -c /home/goatcounter/export-data
