# shared/etc/sv/ntpd/run
#!/bin/sh

exec /usr/sbin/ntpd -n -N -p pool.ntp.org
